# debian/n2view.prerm
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

invoke-rc.d n2view stop || true
#DEBHELPER#

exit 0


