# debian/n2view.postrm
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

[ $1 = purge ] && update-rc.d n2view remove > /dev/null
#DEBHELPER#

exit 0


